# Downloads/farm-admin-enablevnc/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/
environment=DISPLAY=":1",HOME="/root",SHELL="/bin/bash",USER="root",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock

[unix_http_server]
file=/dev/shm/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:xvfb]
priority=1
directory=/tmp
command=/usr/bin/Xvfb :1 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset -nolisten tcp
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err
startsecs=2
startretries=3

[program:x11vnc]
priority=10
command=/usr/bin/x11vnc -display :1 -nopw -forever -shared -rfbport 5900 -listen 0.0.0.0 -xkb -noxrecord -noxfixes -noxdamage -permitfiletransfer
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.err
startsecs=5
environment=DISPLAY=":1",HOME="/root",XAUTHORITY="/root/.Xauthority"

[program:novnc]
priority=15
command=/bin/bash -c "sleep 2 && cd /usr/share/novnc && ./utils/launch.sh --vnc localhost:5900 --listen 0.0.0.0:80 --web /usr/share/novnc"
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.err
startsecs=10
environment=DISPLAY=":1",HOME="/root"

[program:xfce4]
priority=20
command=/usr/bin/startxfce4
environment=DISPLAY=":1",HOME="/root",SHELL="/bin/bash",USER="root",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8",XDG_SESSION_TYPE="x11"
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/xfce4.log
stderr_logfile=/var/log/xfce4.err
startsecs=5

[program:eternalfarm-agent]
priority=50
command=/bin/bash -c "sleep 10 && cd /root && DISPLAY=:1 /usr/local/bin/EternalFarmAgent --auth-agent-key=%(ENV_AUTH_AGENT_KEY)s --show-gui"
autostart=true
autorestart=true
startretries=5
user=root
environment=DISPLAY=":1",
    HOME="/root",
    XAUTHORITY="/root/.Xauthority",
    ETERNAL_API_URL="%(ENV_ETERNAL_API_URL)s",
    AUTH_AGENT_KEY="%(ENV_AUTH_AGENT_KEY)s",
    DREAMBOT_USERNAME="%(ENV_DREAMBOT_USERNAME)s",
    DREAMBOT_PASSWORD="%(ENV_DREAMBOT_PASSWORD)s",
    DREAMBOT_SCRIPT="%(ENV_DREAMBOT_SCRIPT)s",
    DREAMBOT_WORLD="%(ENV_DREAMBOT_WORLD)s",
    DREAMBOT_ARGS="%(ENV_DREAMBOT_ARGS)s",
    JAVA_HOME="/usr/lib/jvm/java-8-openjdk-amd64"
stdout_logfile=/var/log/eternalfarm-agent.log
stderr_logfile=/var/log/eternalfarm-agent.err
startsecs=15

[program:cpu_manager]
priority=60
command=/root/cpu_manager.sh
autostart=true
autorestart=true
startretries=3
user=root
environment=DISPLAY=":1",HOME="/root",XAUTHORITY="/root/.Xauthority"
stdout_logfile=/var/log/cpu_manager.log
stderr_logfile=/var/log/cpu_manager.err
startsecs=20

[program:farm-manager]
priority=70
command=node server.js
directory=/app
user=farmboy
autostart=true
autorestart=true
stderr_logfile=/var/log/farm-manager.err.log
stdout_logfile=/var/log/farm-manager.out.log
environment=NODE_ENV="%(ENV_NODE_ENV)s",
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    API_KEY="%(ENV_API_KEY)s",
    ETERNALFARM_AGENT_KEY="%(ENV_ETERNALFARM_AGENT_KEY)s",
    ETERNAL_API_URL="%(ENV_ETERNAL_API_URL)s",
    DISCORD_WEBHOOK_URL="%(ENV_DISCORD_WEBHOOK_URL)s",
    PORT="%(ENV_PORT)s"
startsecs=15
startretries=3

[program:ssh]
priority=80
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stderr_logfile=/var/log/ssh.err.log
stdout_logfile=/var/log/ssh.out.log

[program:entry-sh]
priority=90
command=/bin/bash -c "sleep 30 && DISPLAY=:1 /root/Entry.sh"
autostart=true
autorestart=false
startretries=1
user=root
environment=DISPLAY=":1"
stdout_logfile=/var/log/entry.log
stderr_logfile=/var/log/entry.err
startsecs=5

[program:eternalfarm-checker]
priority=100
command=/bin/bash -c "sleep 60 && cd /root && xhost + && DISPLAY=:1 /usr/local/bin/EternalFarmChecker --auth-agent-key=%(ENV_AUTH_AGENT_KEY)s --show-gui"
autostart=true
autorestart=true
startretries=5
user=root
environment=DISPLAY=":1",HOME="/root",XAUTHORITY="/root/.Xauthority"
stdout_logfile=/var/log/eternalfarm-checker.log
stderr_logfile=/var/log/eternalfarm-checker.err
startsecs=10

[program:eternalfarm-browser-automator]
priority=110
command=/bin/bash -c "sleep 75 && cd /root && xhost + && DISPLAY=:1 /usr/local/bin/EternalFarmBrowserAutomator --auth-agent-key=%(ENV_AUTH_AGENT_KEY)s --show-gui"
autostart=true
autorestart=true
startretries=5
user=root
environment=DISPLAY=":1",HOME="/root",XAUTHORITY="/root/.Xauthority"
stdout_logfile=/var/log/eternalfarm-browser-automator.log
stderr_logfile=/var/log/eternalfarm-browser-automator.err
startsecs=10

# Core Configuration
NODE_ENV=production
PORT=3001

# DreamBot Settings
DREAMBOT_USERNAME=your_username
DREAMBOT_PASSWORD=your_password
DREAMBOT_SCRIPT=your_script
DREAMBOT_WORLD=301
DREAMBOT_ARGS=--no-splash --developer-mode

# System Optimization
HYPERTHREADS_PER_CLIENT=2
CLIENTS_PER_WORKSPACE=4
CPU_PRIORITY=-10
JAVA_OPTS=-Xmx512m -XX:+UseG1GC

# Display/VNC Settings
RESOLUTION=1920x1080
COLOR_DEPTH=24
VNC_PORT=5900
NOVNC_PORT=8080

# New DreamBot process pattern
DREAMBOT_PROCESS_PATTERN="java.*dreambot"

# New EternalFarm agent pattern
ETERNALFARM_AGENT_PATTERN="EternalFarmAgent"

# Required
ETERNALFARM_AGENT_KEY=your_key
DISCORD_WEBHOOK_URL=your_webhook

# Database
DATABASE_URL=mysql://farmboy:password@mariadb:3306/farm_manager

# Optional but recommended
VNC_PASSWORD=your_password 